# services/clean.timer
[Unit]
Description=Run cleanup every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
